# experiments/230111_twstft_2M5/crontab
06 1-23/2 * * *  /root/go.sh 3
45 1-23/2 * * *  /root/go.sh 3
