# experiments/2403/RPi/crontab
02 0-23/2 * * *  /root/2401/switch_off

